# bbai64_os/recipes/beagleboard/dummy/files/dummy.service
[Unit]
Description=Dummy

[Service]
ExecStart=/usr/bin/dummy.sh
Restart=always

[Install]
WantedBy=multi-user.target
